# arpGuard.service
[Unit]
Description=ARPGuard

[Service]
Type=simple
ExecStart=/opt/arp_guard/ARP_Guard.py

[Install]
WantedBy=multi-user.target
